# deploy/maestro-backup.timer
[Unit]
Description=Hourly MaestroVPN control-plane backup

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
